# deploy/altered-wiki-deploy.service
[Unit]
Description=altered.wiki deploy webhook
After=network.target

[Service]
Type=simple
User=http
WorkingDirectory=/opt/altered-wiki/tools/deploy
ExecStart=/usr/bin/python3 /opt/altered-wiki/tools/deploy/webhook.py
EnvironmentFile=/etc/altered-wiki/deploy.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
